# downloadPictures.timer
[Unit]
Description=Download new pictures daily and at boot

[Timer]
OnCalendar=*-*-* 9:00
OnBootSec=60
RandomizedDelaySec=500
